# notebook.conf
[supervisord]

[program:notebook]
command=/home/vagrant/.pyenv/shims/ipython3 notebook --ip 0.0.0.0 --no-browser
autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=5MB
user=vagrant
